# systemd/audit-daemon.service
# systemd/audit-daemon.service
[Unit]
Description=Audit Command Tracking Daemon (dispatcher mode)
After=auditd.service
PartOf=auditd.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/audit-daemon /etc/audit-daemon/audit-daemon.conf
Restart=no
User=root
StandardInput=null
StandardOutput=journal
StandardError=journal
KillMode=process
