# systemd/audit-daemon.service
// systemd/audit-daemon.service
[Unit]
Description=Audit Command Tracking Daemon (standalone mode)
After=network.target
ConditionPathExists=!/etc/audisp/plugins.d/audit-daemon.conf

[Service]
Type=simple
ExecStart=/usr/local/sbin/audit-daemon /etc/audit-daemon/audit-daemon.conf
Restart=always
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal
KillMode=process

[Install]
WantedBy=multi-user.target
